# extra/foreman-proxy-abrt-send.cron
# Send collected ABRT reports once every 30 minutes
*/30 * * * * foreman-proxy smart-proxy-abrt-send >>/var/log/foreman-proxy/cron.log 2>&1
