# extra/foreman-proxy-abrt-send.cron
# Send collected ABRT reports once every 10 minutes
*/30 * * * * foreman-proxy smart-proxy-abrt-send
